# production/upstart/805app.conf
#/etc/init/805app.conf
description "My 805 Node Server"
author "steve@mudio.org"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

script
    cd /data/html/805/bootcamp/
    exec sudo -u nodeuser NODE_ENV=production /usr/bin/node /data/html/805/bootcamp/app.js 1> /home/nodeuser/log/805_access.log 2> /home/nodeuser/log/805_error.log
end script
