# upstart/805server.conf
#/etc/init/myapp.conf
description "805 Node.js"
author "steve@bigfrog.net"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

script
        cd /home/ubuntu/myapp
        exec sudo -u ubuntu NODE_ENV=production /usr/bin/node /home/ubuntu/myapp/server.js >> /home/ubuntu/myapp/log/server.log 2>&1
end script
